# aeroserver.service
# Fichier de service systemd pour AeroServer
#
# Installation :
#   1. Ajuster SERVER_DIR ci-dessous avec le chemin absolu du serveur sur Linux
#   2. Créer l'utilisateur système si nécessaire :
#        sudo useradd -r -m -d /opt/aeroserver -s /bin/bash minecraft
#   3. Copier ce fichier :
#        sudo cp aeroserver.service /etc/systemd/system/
#   4. Recharger systemd et activer le service :
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now aeroserver

[Unit]
Description=AeroServer — Serveur Minecraft NeoForge 1.21.1
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minecraft
Group=minecraft

# Adapter ce chemin au répertoire du serveur sur la machine Linux
WorkingDirectory=/opt/aeroserver
ExecStart=/opt/aeroserver/start.sh

# Redémarrer automatiquement en cas de crash (pas si arrêt propre)
Restart=on-failure
RestartSec=15

# Laisser le temps au serveur de s'arrêter proprement avant SIGKILL
TimeoutStopSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aeroserver

[Install]
WantedBy=multi-user.target
